# tests/example_startup_script/etc/systemd/system/example.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/start-example-service.sh

[Install]
WantedBy=default.target
